# src/import_cronjob/pegelwacht.cron
*/16 * * * * /home/pegelwacht/git/pegelwacht/cronjob_scripts/import.sh